# tracker/dpkg/sbc-tracker-0.5.1/etc/systemd/system/ft8-tracker.service
[Unit]
Description=FT8 Mobile Tracker
After=network.target gpsd.service
Wants=gpsd.service

[Service]
Type=simple
User=sbc-tracker
WorkingDirectory=/opt/sbc-tracker
ExecStart=/opt/sbc-tracker/venv/bin/python3 /opt/sbc-tracker/src/main.py -c /etc/sbc-tracker/tracker.conf
Restart=always
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ft8-tracker

[Install]
WantedBy=multi-user.target
